# opencode/llama-cpp.service
[Unit]
Description=Llama.cpp Server

[Service]
Type=simple
WorkingDirectory=%h/.llama/llama.cpp/cuda-12.8
ExecStart=/bin/bash -c "export LD_LIBRARY_PATH=%h/.llama/llama.cpp/cuda-12.8:/usr/local/lib/ollama/cuda_v12 && exec %h/.llama/llama.cpp/cuda-12.8/llama-server -m %h/.llama/models/Qwen3-14B-Q4_K_M.gguf --host 0.0.0.0 --port 8080 -c 131072 -ngl 999 --temp 1.0 --top-p 0.95 --top-k 40 --min-p 0.01"
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
